# Plugin/kurovpn-bot.service
[Unit]
Description=KUROVPN Telegram Bot Controller
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/kurovpn
ExecStart=/usr/bin/python3 /opt/kurovpn/bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
